# units/auto-install.service
[Unit]
Description=Auto-run Arch ARM installer at first boot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/root/auto-install.sh
RemainAfterExit=yes
ConditionPathExists=!/root/.installer_done

[[Install]
